# scenarios.json
{
  "T1059.001": {
    "name": "PowerShell Execution",
    "description": "Adversaries may abuse PowerShell to execute various commands and scripts. PowerShell is a powerful interactive command-line shell and scripting language that is built on the .NET framework. Adversaries use it for a wide range of activities, including discovery, execution, and persistence, often leveraging its ability to run code in memory to bypass traditional disk-based detection.",
    "data_sources": [
      "Process monitoring (Event ID 4688)",
      "PowerShell Script Block Logging (Event ID 4104)",
      "PowerShell Transcription (Event ID 4103)",
      "Command-line arguments"
    ],
    "next_stages": [
      {
        "stage": "Persistence - Scheduled Task",
        "technique": "T1053.005",
        "query": "process.name: \"schtasks.exe\" OR process.command_line: \"/create\""
      },
      {
        "stage": "Credential Access - LSASS Dump",
        "technique": "T1003.001",
        "query": "process.name: \"mimikatz.exe\" OR process.target.name: \"lsass.exe\""
      }
    ]
  },
  "T1053.005": {
    "name": "Scheduled Task",
    "description": "Adversaries may abuse the Windows Task Scheduler to facilitate initial or recurring execution of malicious code. By creating a scheduled task, an attacker can ensure their payload runs at specific intervals, upon system boot, or when a specific user logs in. This provides a reliable method for maintaining persistence and can also be used for remote execution if the attacker has administrative privileges.",
    "data_sources": [
      "Scheduled Task creation events (Event ID 4698)",
      "Process monitoring for schtasks.exe",
      "Windows Registry modifications (HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tree)",
      "File system monitoring for task files in C:\\Windows\\System32\\Tasks"
    ],
    "next_stages": [
      {
        "stage": "Execution - Command and Scripting Interpreter",
        "technique": "T1059",
        "query": "process.name: \"cmd.exe\" OR process.name: \"powershell.exe\""
      }
    ]
  },
  "T1003.001": {
    "name": "LSASS Memory",
    "description": "Adversaries may attempt to access credential material stored in the process memory of the Local Security Authority Subsystem Service (LSASS). This process often contains plaintext passwords, NTLM hashes, and Kerberos tickets for logged-in users. Attackers typically use tools like Mimikatz or built-in Windows utilities like Procdump to create a memory dump of LSASS, which is then parsed to extract sensitive credentials.",
    "data_sources": [
      "Process access events (Event ID 10 in Sysmon)",
      "Process creation events for mimikatz.exe or procdump.exe",
      "File system monitoring for .dmp files created by LSASS",
      "Handle monitoring for lsass.exe"
    ],
    "next_stages": [
      {
        "stage": "Exfiltration - Over Alternative Protocol",
        "technique": "T1048",
        "query": "network.protocol: \"ftp\" OR network.protocol: \"http\""
      }
    ]
  }
}
